# infra/php/supervisord.conf
[supervisord]
nodaemon=true

[program:send-notification-workers-watcher]
command = php /home/dockeruser/backend/daemon/daemon_worker_watcher.php
autorestart = true
priority = 1